# sv/log/run
#!/bin/sh
#
#

set -u -e
dir=/var/log
sv_name=miniuni_http

mkdir -p         $dir/$sv_name
exec svlogd -ttt $dir/$sv_name
